# fedora-42-general/config/systemd/user/sunset-theme.service
[Unit]
Description=Sonnenstand prüfen und GTK color-scheme (Firefox etc.) umschalten

[Service]
Type=oneshot
ExecStart=%h/.local/bin/sunset-theme
